# root/etc/s6-overlay/s6-rc.d/svc-httpserver/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc \
  cd /data/tunasync/mirror python3 -m http.server -b :: ${HTTP_PORT}
