# root/etc/s6-overlay/s6-rc.d/svc-tunasync-manager/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc \
  tunasync manager --config /config/tunasync/manager.conf
